# HomeLab/roles/sudo/templates/sudo-io-replay-rotate.timer.j2
[Unit]
Description=Sudo IO Rotate

[Timer]
OnCalendar=daily
RandomizedDelaySec={{ sudo_random_delay_time }}
Persistent=true

[Install]
WantedBy=timers.target
