# HomeLab/roles/aide/files/aide-check.service
[Unit]
Description=AIDE check
[Service]
Type=oneshot
Slice=aide-check.slice
ExecStart=/usr/sbin/aide --check
