# HomeLab/roles/aide/templates/aide-check.timer.j2
[Unit]
Description=AIDE check

[Timer]
OnCalendar=weekly
RandomizedDelaySec={{ aide_random_delay_time }}
FixedRandomDelay=true
Persistent=true

[Install]
WantedBy=timers.target
